# proxy/cc-proxy.service.in
[Unit]
Description=Clear Containers Proxy

[Service]
ExecStart=@libexecdir@/cc-proxy

[Install]
WantedBy=multi-user.target
